# 02_build/security/hardening/docker-user-firewall.service
[Unit]
Description=Docker-User Firewall Rules (AMP-289)
Documentation=https://linear.app/amplifiedpartners/issue/AMP-289
# Run after Docker (it creates the DOCKER-USER chain on start).
# PartOf ensures rules re-apply whenever Docker restarts.
After=docker.service
PartOf=docker.service

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/opt/amplified/security/docker-user-firewall.sh
ExecReload=/opt/amplified/security/docker-user-firewall.sh

[Install]
WantedBy=multi-user.target docker.service

# Signed-by: Devon-abaf | 2026-05-09 | session devin-abaf337e1eda492f85af1d26fa1614b9
